# build/overlay/etc/systemd/system/jetson-oc-perf.service
[Unit]
Description=Jetson Nano Performance Tuning
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo performance > /sys/devices/system/cpu/cpufreq/policy0/scaling_governor && echo 2014500 > /sys/devices/system/cpu/cpufreq/policy0/scaling_max_freq && echo 1152000000 > /sys/devices/57000000.gpu/devfreq/57000000.gpu/max_freq && echo noop > /sys/block/mmcblk0/queue/scheduler'

[Install]
WantedBy=multi-user.target
